# dotfiles/etc/cron.hourly/btrfs-snap
#!/bin/sh
/usr/local/sbin/btrfs-snap /home hourly 12
