# dotfiles/etc/cron.hourly/consul-leave
#!/bin/sh
exec chpst -u consul /usr/local/bin/consul leave >/dev/null
