# dotfiles/etc/cron.hourly/ntpdate
#!/bin/sh
ntpdate -sB hu.pool.ntp.org europe.pool.ntp.org pool.ntp.org
