# dotfiles/etc/cron.daily/btrfs-snap
#!/bin/sh
/usr/local/sbin/btrfs-snap /home daily 5
