# dotfiles/etc/service/godoc/run
#!/bin/sh
env GOMAXPROCS=${GOMAXPROCS:-2} GOPATH=${GOPATH:-$HOME} godoc -http=:6060 
