# dotfiles/etc/service/godoc/run
#!/bin/sh
env GOPATH=${GOPATH:-$HOME} godoc -http=:6060 -analysis=type,pointer
