# dotfiles/etc/systemd/system/runsvdir.service
[Unit]
Description=runit-runsvdir

[Service]
Type=simple
ExecStart=/usr/sbin/runsvdir-start
Restart=always

[Install]
WantedBy=basic.service
