# dotfiles/etc/cron.weekly/btrfs-snap
#!/bin/sh
/usr/local/sbin/btrfs-snap /home weekly 3
