# dotfiles/etc/cron.weekly/btrfs-scrub
#!/bin/sh
awk '/btrfs/ {print $2}' /proc/mounts | while read dn; do
  echo $dn
  btrfs scrub start -B $dn
done
